# dist/go_audit_container_debian/DEBIAN/postinst
#!/bin/bash

systemctl disable auditd
systemctl stop auditd
systemctl daemon-reload
systemctl enable goaudit
systemctl start goaudit
